# system/ainod.socket
[Socket]
ListenStream=@ainod

[Install]
WantedBy=sockets.target
